# base/etc/services.d/openbox/run
#!/usr/bin/with-contenv sh
exec \
	s6-setuidgid container wait-for localhost:5900 -t 5 -- /usr/bin/openbox
